# ansible/roles/plex_server_config/templates/plex-config-sync.service.j2
# ==============================================================================
# File: /etc/systemd/system/plex-config-sync.service
# Purpose:
#   Sync Plex server configuration from the primary Plex server.
# Notes:
#   - Generated by Ansible role plex_server_config.
#   - Plex does not provide official multi-server config replication.
# ==============================================================================

[Unit]
Description=Sync Plex configuration from {{ plex_server_config_sync_primary }}
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ plex_server_config_sync_user }}
ExecStart=/bin/bash -lc 'rsync {{ plex_server_config_sync_rsync_options | join(" ") }} {% for item in plex_server_config_sync_excludes %}--exclude={{ item | quote }} {% endfor %}{{ plex_server_config_sync_user }}@{{ hostvars[plex_server_config_sync_primary].ansible_host | default(plex_server_config_sync_primary) | quote }}:{{ plex_server_config_sync_remote_path | quote }} {{ plex_server_config_sync_local_path | quote }}'
ExecStartPost=/bin/systemctl try-restart {{ plex_server_config_service_name }}.service
