# ansible/roles/plex_server_config/templates/plex-config-sync.timer.j2
# ==============================================================================
# File: /etc/systemd/system/plex-config-sync.timer
# Purpose:
#   Timer for Plex configuration sync from the primary Plex server.
# Notes:
#   - Generated by Ansible role plex_server_config.
# ==============================================================================

[Unit]
Description=Run Plex configuration sync on schedule

[Timer]
OnCalendar={{ plex_server_config_sync_schedule }}
Persistent=true
Unit=plex-config-sync.service

[Install]
WantedBy=timers.target
